# fleetscorebot/fleetscore@staging.timer
[Unit]
Description=Run fleetscore for staging every hour

[Timer]
OnCalendar=*-*-* *:00:00
AccuracySec=10s

[Install]
WantedBy=timers.target
