# fleetscorebot/fleetscore@prod.timer
[Unit]
Description=Run fleetscore for production every hour

[Timer]
OnCalendar=*-*-* *:05:00
AccuracySec=10s

[Install]
WantedBy=timers.target
